# crontab
0 * * * * curl http://raja-nournia.dotcloud.com/update >> cron.log